# retired/process_batches.cron
#!/bin/bash -l
SCRIPT_DIR=$( cd -- "$( dirname -- "${BASH_SOURCE[0]}" )" &> /dev/null && pwd )
LOCKFILE=$SCRIPT_DIR/processing.lock

export PATH=$SCIRPT_DIR:$HOME/bin:$PATH

if [ -e $LOCKFILE ]; then
    exit 0
fi

echo $$ > $LOCKFILE

$SCRIPT_DIR/process_batches.py $@
RC=$?

rm $LOCKFILE

exit $RC
